# automation/io.kern.listen-on-camera.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>io.kern.listen-on-camera</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/akern/.dotfiles/automation/listen-for-camera.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>ELGATO_LIGHT_HOST</key>
      <string>192.168.0.199</string>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/listen-on-camera.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/listen-on-camera.error.log</string>
  </dict>
</plist>
